# crontab
* * * * * cd /app && /root/.bun/bin/bun scripts/cron/expire-invitations.ts >> /var/log/cron.log 2>&1
