# crontab
* * * * * cd /app && bun scripts/cron/expire-invitations.ts >> /var/log/cron.log 2>&1
